# iso-builder/overlays/et-r5-general/overlay/etc/skel/.config/autostart/udiskie.desktop
[Desktop Entry]
Type=Application
Name=USB Automount
Exec=udiskie -a -n -t
Hidden=false
NoDisplay=true
X-GNOME-Autostart-enabled=true
