# iso-builder/overlays/et-r5-general/overlay/lib/systemd/system/rigctld.service
[Unit]
Description=Rig Control Service
After=multi-user.target

[Service]
Type=simple
ExecStart=/opt/emcomm-tools/sbin/wrapper-rigctld.sh start

[Install]
WantedBy=multi-user.target
